# nnn.conf
[program:nnn]
command=python3 /home/curator2/nnn/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/nnn/test.err.log
stdout_logfile=/home/curator2/nnn/test.out.log